# init/deluge-web
#!/sbin/openrc-run
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/openssh/files/sshd.rc6.4,v 1.5 2015/05/04 02:56:25 vapier Exp $

description="Deluge Web interface"

start() {

	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --background --start \
	    --exec "/usr/bin/deluge-web" \
	    --user deluge \
	    --make-pidfile --pidfile /var/run/deluge-web.pid \
	    -- ${DELUGEWEB_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --quiet --stop --exec "/usr/bin/deluge-web" \
	    --retry=TERM/30/KILL/5 \
	    --pidfile /var/run/deluge-web.pid
	eend $?
}
