# init/openvpn
#!/sbin/openrc-run
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="VPN service"

vpn_provider="$(echo $OPENVPN_PROVIDER | tr '[A-Z]' '[a-z]')"
vpn_provider_configs="/etc/openvpn/$vpn_provider"
vpn_config_file=""

if [ -n "${vpn_provider}" ] && [ ${SVCNAME} != "openvpn" ]; then
	vpnpid="/var/run/openvpn.${vpn_provider}.pid"
else
	vpnpid="/var/run/openvpn.pid"
fi

depend() {
	need localmount net
	use dns
	after bootmisc
}

checkconfig() {

	if [ ! -c "/dev/net/tun" ]; then
		if [ ! -d "/dev/net" ]; then
			mkdir -p /dev/net
		fi
		mknod /dev/net/tun c 10 200
	fi

	if [ ! -c "/dev/net/tun" ]; then
		eerror "TUN interface not available."
		return 1
	fi

	if [ ! -d "$vpn_provider_configs" ]; then
		eerror "Could not find OpenVPN provider: $OPENVPN_PROVIDER"
		return 1
	fi

	if [ ! -z "$OPENVPN_CONFIG" ]; then
		if [ -f $vpn_provider_configs/"${OPENVPN_CONFIG}".ovpn ]; then
			vpn_config_file=$vpn_provider_configs/${OPENVPN_CONFIG}.ovpn
		else
			vpn_config_file=$vpn_provider_configs/default.ovpn
		fi
	else
		vpn_config_file=$vpn_provider_configs/default.ovpn
	fi

	if [ "${OPENVPN_USERNAME}" = "**None**" ] || [ "${OPENVPN_PASSWORD}" = "**None**" ] ; then
		eerror "OpenVPN credentials not set. Exiting."
		return 1
	else
		mkdir -p /config
		echo $OPENVPN_USERNAME > /config/openvpn-credentials.txt
		echo $OPENVPN_PASSWORD >> /config/openvpn-credentials.txt
		chmod 600 /config/openvpn-credentials.txt
	fi
	
	if [ -z ${vpn_config_file} ] || [ ! -f "${vpn_config_file}" ]; then
		eerror "VPN Configuration $vpn_config_file does" \
			"not exist."
		return 1
	fi
	
  	if [ ! -z "${PUID}" ]; then                                                                                                                                                                 
    		if [ ! "$(id -u deluge)" -eq "${PUID}" ]; then
    
      			# usermod likes to chown the home directory, so create a new one and use that
      			# However, if the new UID is 0, we can't set the home dir back because the
      			# UID of 0 is already in use (executing this script).
      
      			if [ ! "${PUID}" -eq 0 ]; then
        			mkdir /tmp/temphome
        			usermod -d /tmp/temphome deluge
      			fi
      
      			# Change the UID
      			usermod -o -u "${PUID}" deluge
      
      			# Cleanup the temp home dir
      			if [ ! "${PUID}" -eq 0 ]; then
        			usermod -d /config deluge
        			rm -Rf /tmp/temphome
      			fi
			
    		fi
  	fi

  	if [ ! -z "${PGID}" ]; then
    		if [ ! "$(id -g deluge)" -eq "${PGID}" ]; then
      			groupmod -o -g "${PGID}" users
    		fi
  	fi

	chown deluge:users /data

}

start() {
	# If we are re-called by the openvpn gentoo-up.sh script
	# then we don't actually want to start openvpn
	if [ "${IN_BACKGROUND}" = "true" ]; then
		return 0
	fi
	
	ebegin "Starting ${SVCNAME}"

	checkconfig || return 1

	local args="" reenter=${RE_ENTER:-no}
	# If the config file does not specify the cd option, we do
	# But if we specify it, we override the config option which we do not want

	if ! grep -q "^[ 	]*cd[ 	].*" "${vpn_provider_configs}" ; then
		config_path=$(dirname "$vpn_provider_configs")
		args="${args} --cd ${config_path}"
	fi

	# We mark the service as inactive and then start it.
	# When we get an authenticated packet from the peer then we run our script
	# which configures our DNS if any and marks us as up.

	reenter="yes"
	args="${args} --up-delay --up-restart"
	args="${args} --script-security 2"
	args="${args} --up /etc/openvpn/up.sh"
	args="${args} --down-pre --down /etc/openvpn/down.sh"

	# Warn about the inability to change ip/route/dns information when
	# dropping privs
	if grep -q "^[ 	]*user[ 	].*" "${vpn_config_file}" ; then
		ewarn "WARNING: You are dropping root privileges!"
		ewarn "As such openvpn may not be able to change ip, routing"
		ewarn "or DNS configuration."
	fi

	# Ensure that our scripts get the PEER_DNS variable
	[ -n "${PEER_DNS}" ] && args="${args} --setenv PEER_DNS ${PEER_DNS}"

	[ "${reenter}" = "yes" ] && mark_service_inactive "${SVCNAME}"

	#oma
	if [ -n "${LOCAL_NETWORK-}" ]; then
		eval $(/sbin/ip r l m 0.0.0.0 | awk '{if($5!="tun0"){print "GW="$3"\nINT="$5; exit}}')
		if [ -n "${GW-}" -a -n "${INT-}" ]; then
			/sbin/ip r a "$LOCAL_NETWORK" via "$GW" dev "$INT"
		fi
	fi

	start-stop-daemon --start --exec /usr/sbin/openvpn --pidfile "${vpnpid}" \
		-- --config "${vpn_config_file}" --writepid "${vpnpid}" --daemon \
		--setenv SVCNAME "${SVCNAME}" ${args}
	eend $? "Check your logs to see why startup failed"
}

stop() {
	# If we are re-called by the openvpn gentoo-down.sh script
	# then we don't actually want to stop openvpn
	if [ "${IN_BACKGROUND}" = "true" ] ; then
		mark_service_inactive "${SVCNAME}"
		return 0
	fi

	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --quiet \
		--exec /usr/sbin/openvpn --pidfile "${vpnpid}"
	eend $?
}

# vim: set ts=4 :
